# cron.d/backup_mysql_docker
00 3 * * *  root bash /opt/acme-php-nginx-dockerized/backup-db.sh > /dev/null 2>&1
