# cron.d/backup_mysql_docker
00 3 * * *  root bash /root/backup-db.sh > /dev/null 2>&1
